# crontab
# Fix files that are set with invalid ownership for cron.
*/5 * * * * root (chown root: /etc/cron.d/* && chmod 0600 /etc/cron.d/*) >/dev/null 2>&1
